# systemd/inside-qa.service
[Unit]
Description=Inside Pipeline — QA Agent (Codex CLI: Review + Test + BRD check)
Documentation=https://github.com/ppxtestarena4/inside
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/home/agent/<repo>

ExecStart=/home/agent/<repo>/pipeline/qa-daemon.sh

# Переменные окружения
Environment=OPENAI_API_KEY=
Environment=PIPELINE_REPO=ppxtestarena4/inside
Environment=SLEEP_INTERVAL=300
Environment=HOME=/home/agent
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/agent/.local/bin

Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=inside-qa

TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
